# systemd/shellweb@.service
[Service]
ExecStart=/home/jianjun/code/shellweb/http.sh
User=jianjun
Group=jianjun
StandardInput=socket
StandardOutput=socket
