# socket/shellweb@.service
[Service]
ExecStart=/root/code/shellweb/http.sh
User=featen
Group=featen
StandardInput=socket
StandardOutput=socket
